# target/etc/supervisor/conf.d/named.conf
[program:named]
command=named -f
redirect_stderr=true
numprocs=1
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
